# 71-ft9201_init.rules
ACTION=="add", SUBSYSTEMS=="usb", DRIVERS=="usb", ATTRS{idVendor}=="2808", ATTRS{idProduct}=="93a9", RUN+="/usr/bin/ft9201_init.py"
